# dataCreator/crontab
* * * * * cd /project && /usr/local/bin/php artisan schedule:run >> /project/cron.log 2>&1

